# keesync.service
[Unit]
Description=keesync - keepassx db synchronization
After=network.target

[Service]
User={USER}
Group={USER}
WorkingDirectory={WORKING_DIRECTORY}
Environment="KEESYNC_APP_KEY={APP_KEY}"
ExecStart={PYTHON_PATH} {WORKING_DIRECTORY}/keesync.py -p {DATABASE_PATH} -l info -s 1
Restart=on-failure

[Install]
WantedBy=multi-user.target
